# sway/.config/systemd/user/wob.socket.d/override.conf
[Socket]
RemoveOnStop=on
# If wob crashes from an input systemd should NOT shove it right back into it after it restarted
FlushPending=yes
