# borg/.config/systemd/user/automatic-backup@.timer
[Unit]
Description=Run automatic backup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=graphical-session.target

