# borg/.config/systemd/user/automatic-backup@.service
[Unit]
Description=Borgmatic backup for %i
Before=shutdown.target

[Service]
ExecStart=/home/scrumplex/bin/automatic-backup %i
NotifyAccess=all
