# sway/.config/systemd/user/swayidle.service
[Unit]
Description=Sway Idle Daemon

[Service]
ExecStart=/usr/bin/swayidle

[Install]
WantedBy=sway-session.target

