# mpd/.config/systemd/user/mpdscribble.service.d/override.conf
[Install]
WantedBy=
WantedBy=default.target