# mpv/.config/systemd/user/jellyfin-mpv-shim.service
[Unit]
Description=Jellyfin MPV Shim

[Service]
ExecStart=/usr/bin/jellyfin-mpv-shim

[Install]
WantedBy=graphical-session.target

